# deploy/shurli-daemon.service
[Unit]
Description=Shurli daemon - P2P network service
Documentation=https://shurli.io
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/shurli daemon
Restart=on-failure
RestartSec=5
WatchdogSec=90
User=shurli
Group=shurli

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/shurli /run/user
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
