# relay-server/relay-server.service
[Unit]
Description=peer-up Relay Server
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
WatchdogSec=90

# Run as the user who owns the repo - NEVER as root
# Change this to your username:
User=YOUR_USERNAME
Group=YOUR_USERNAME

# Set WorkingDirectory to where relay-server.yaml lives
# Change this to your actual path:
WorkingDirectory=/home/YOUR_USERNAME/peer-up/relay-server

# Single binary - relay mode
ExecStart=/home/YOUR_USERNAME/peer-up/relay-server/peerup relay serve

Restart=always
RestartSec=5

# Increase file descriptor limit for many concurrent connections
LimitNOFILE=65536

# Allow binding to privileged ports (443) without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Security hardening - restrict what the service can do
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow read/write only to the working directory (for key file, config, logs)
ReadWritePaths=/home/YOUR_USERNAME/peer-up/relay-server

# Logging via journald (configure rotation in /etc/systemd/journald.conf)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=peerup-relay

[Install]
WantedBy=multi-user.target
